# voice-trimmer-bot.service
[Unit]
Description=Voice Trimmer Telegram Bot
Documentation=https://github.com/mohamadbastaniyan/voiceTrimmer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/voiceTrimmer
Environment="PATH=/root/voiceTrimmer/venv/bin"
Environment="TELEGRAM_BOT_TOKEN=YOUR_TOKEN_HERE"
Environment="MAX_DURATION=60"
Environment="LOG_DIR=/root/voiceTrimmer/logs"
Environment="LOG_LEVEL=INFO"
ExecStart=/root/voiceTrimmer/venv/bin/python /root/voiceTrimmer/src/bot.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-trimmer-bot

NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/root/voiceTrimmer/logs /tmp/voice_trimmer

[Install]
WantedBy=multi-user.target
